# trading-bot.service
[Unit]
Description=Crypto Trading Bot
After=network.target postgresql.service
Wants=postgresql.service
Documentation=https://github.com/your-repo/TradingBot

[Service]
Type=simple
User=elcrypto
Group=elcrypto
WorkingDirectory=/home/elcrypto/TradingBot

# Environment
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/elcrypto/TradingBot/.env

# Execute bot using venv python
ExecStart=/home/elcrypto/TradingBot/venv/bin/python /home/elcrypto/TradingBot/main.py --mode production

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
MemoryLimit=1G
CPUQuota=50%

# Logging
StandardOutput=append:/home/elcrypto/TradingBot/logs/systemd-output.log
StandardError=append:/home/elcrypto/TradingBot/logs/systemd-error.log
SyslogIdentifier=trading-bot

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
